# nx-workspace/apps/backend/pfbackend/supervisord.conf
[unix_http_server]
file=/tmp/supervisor.sock   ; the path to the socket file

[supervisord]
logfile=/tmp/supervisord.log ; main log file; default $CWD/supervisord.log
loglevel=info                ; log level; default info; others: debug,warn,trace
pidfile=/tmp/supervisord.pid ; supervisord pidfile; default supervisord.pid

[rpcinterface:supervisor]
supervisor.rpcinterface_factory = supervisor.rpcinterface:make_main_rpcinterface

[supervisorctl]
serverurl=unix:///tmp/supervisor.sock ; use a unix:// URL  for a unix socket

[program:gunicorn]
command=/home/joeym/Repos/pfcharacter/nx-workspace/apps/backend/env/bin/gunicorn pfcharacter.wsgi
directory=/home/joeym/Repos/pfcharacter/nx-workspace/apps/pfcharacter
user=joey_malme
autostart=true
autorestart=true
redirect_stderr=true